# scripts/com.rui.notion-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.rui.notion-backup</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/rui.huang/Documents/Claude/Projects/Notion_Auto_PM/scripts/run_backup.sh</string>
    </array>

    <!-- Run daily at 6:00 AM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- If the machine was asleep at 6am, run when it next wakes -->
    <key>StartInterval</key>
    <integer>86400</integer>

    <key>StandardOutPath</key>
    <string>/Users/rui.huang/Documents/Claude/Projects/Notion_Auto_PM/backups/launchd_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/rui.huang/Documents/Claude/Projects/Notion_Auto_PM/backups/launchd_stderr.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
